# homeserver/systemd-containers/nextcloud/nextcloud-proxy.container
# nextcloud-proxy.container
[Unit]
Description=nextcloud
Requires=nextcloud.service
After=nextcloud.service

[Container]
Image=docker.io/library/caddy:latest
AutoUpdate=registry
Pod=nextcloud.pod
Volume=./caddy/Caddyfile:/etc/caddy/Caddyfile:Z
Volume=./caddy/cert.pem:/etc/ssl/certs/cert.pem:Z
Volume=./caddy/cert.key:/etc/ssl/private/cert.key:Z
LogDriver=none

[Service]
Restart=always
